# deploy/config/systemd/rag-ui.service.template
"""
Systemd service configuration for RAG Support Client UI (Streamlit)
Includes security hardening, resource limits and restart policies.

Author: Oliver Marshall
Date: December 15, 2024
"""

[Unit]
Description=RAG Support Client Streamlit UI Service
Documentation=https://github.com/espritdunet/rag-support-client
After=network.target network-online.target rag-api.service
Requires=network-online.target rag-api.service
# If using remote Ollama, remove the following line
After=ollama.service
Requires=ollama.service

[Service]
# Service execution
Type=simple
User=${RAG_USER}
Group=${RAG_GROUP}
WorkingDirectory=${INSTALL_BASE}/current
Environment="PATH=${INSTALL_BASE}/current/.venv/bin"
EnvironmentFile=${CONFIG_BASE}/.env
Environment="STREAMLIT_BROWSER_GATHER_USAGE_STATS=false"
ExecStart=${INSTALL_BASE}/current/.venv/bin/streamlit run rag_support_client/streamlit/app.py \
    --server.port ${STREAMLIT_PORT} \
    --server.address ${STREAMLIT_HOST} \
    --server.maxUploadSize 5 \
    --server.maxMessageSize 200 \
    --server.enableXsrfProtection true \
    --server.enableCORS false \
    --browser.serverAddress ${NGINX_SERVER_NAME} \
    --browser.gatherUsageStats false \
    --logger.level info \
    --theme.base light

# Restart configuration
Restart=always
RestartSec=5
StartLimitInterval=0
StartLimitBurst=5

# Resource management
CPUQuota=100%
CPUWeight=80
MemoryMax=2G
TasksMax=2048
LimitNOFILE=32768
LimitNPROC=2048
IOWeight=80

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Logging
StandardOutput=append:${LOG_BASE}/ui.log
StandardError=append:${LOG_BASE}/ui.error.log
SyslogIdentifier=rag-ui

# Performance tuning
Nice=0
IOSchedulingClass=best-effort
IOSchedulingPriority=4

# Watchdog configuration
WatchdogSec=60s
TimeoutStartSec=60s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
